# etc/rc.d/automountd
#!/bin/ksh

daemon="/usr/sbin/automountd"

. /etc/rc.d/rc.subr

rc_post() {
	rm -f /var/run/automountd.pid
}

rc_cmd $1
