# etc/rc.d/autounmountd
#!/bin/ksh

daemon="/usr/sbin/autounmountd"

. /etc/rc.d/rc.subr

rc_post() {
	rm -f /var/run/autounmountd.pid
}

rc_cmd $1
